# .config/autostart/ibus.desktop
[Desktop Entry]
Name=IBus
Exec=ibus-daemon -drx --panel=/usr/lib/kimpanel-ibus-panel
Terminal=false
Type=Application
StartupNotify=false
